# metrics-exporter.service
[Unit]
Description=Platform-Agnostic Metrics Exporter
After=network.target

[Service]
Type=simple
User=otelcol
Group=otelcol
WorkingDirectory=/opt/metrics-exporters
ExecStart=/opt/metrics-exporters/venv/bin/python /opt/metrics-exporters/main.py
Restart=on-failure
RestartSec=5
Environment=PYTHONPATH=/opt/metrics-exporters

# Export format configuration (choose one):
# For OTLP export (default):
Environment=EXPORT_FORMAT=otlp
Environment=OTEL_ENDPOINT=infra-metrics1.home.seiffert.me:4317
Environment=OTEL_SERVICE_NAME=metrics-exporter
Environment=OTEL_SERVICE_VERSION=1.0.0
Environment=OTEL_INSECURE=true

# For Prometheus export (uncomment and configure):
#Environment=EXPORT_FORMAT=prometheus

# Collector configuration (optional):
#Environment=ENABLED_COLLECTORS=memory,disk,process,cpu,network
#Environment=CPU_COLLECTION_METHOD=auto
#Environment=CPU_MEASUREMENT_INTERVAL=15.0
#Environment=NETWORK_INTERFACES=eth0,ens3
#Environment=NETWORK_EXCLUDE_INTERFACES=lo,docker0
#Environment=LXC_DETECT_CONTAINER_LIMITS=true
#Environment=LXC_PREFER_CGROUP_STATS=true
#Environment=MAX_MEASUREMENT_HISTORY=10

[Install]
WantedBy=multi-user.target