# metrics-exporter.service
[Unit]
Description=Metrics Exporter (OTLP Only)
After=network.target

[Service]
Type=simple
User=otelcol
Group=otelcol
WorkingDirectory=/opt/metrics-exporters
ExecStart=/opt/metrics-exporters/venv/bin/python /opt/metrics-exporters/main.py
Restart=on-failure
RestartSec=5
Environment=PYTHONPATH=/opt/metrics-exporters

# OTLP Configuration (required)
Environment=OTLP_ENDPOINT=infra-metrics1.home.seiffert.me:4317
Environment=SERVICE_NAME=metrics-exporter
Environment=SERVICE_VERSION=1.0.0
Environment=OTLP_INSECURE=true

# Optional settings
Environment=COLLECTION_INTERVAL=30
Environment=ENABLED_COLLECTORS=memory,cpu,filesystem,network,process
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target